# scripts/oiko-deploy.service
[Unit]
Description=Oiko Auto-Deploy Service
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
ExecStart=/opt/oiko/scripts/auto-deploy.sh
StandardOutput=journal
User=root
