# scripts/oiko-deploy.timer
[Unit]
Description=Oiko Auto-Deploy Timer (check every 5 minutes)
After=network-online.target docker.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=oiko-deploy.service

[Install]
WantedBy=timers.target
